# launchd/com.earworm.autogen.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Auto-topic generator: proposes 3 fresh topics every Monday at 07:00.
     Caps output at ~3 episodes/week (<=30 min) given the weekday drain. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.earworm.autogen</string>
    <key>ProgramArguments</key>
    <array>
        <string>__EARWORM_DIR__/.venv/bin/earworm</string>
        <string>autogen</string>
        <string>--count</string>
        <string>3</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__EARWORM_DIR__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>1</integer>
        <key>Hour</key><integer>7</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__EARWORM_DIR__/logs/autogen.log</string>
    <key>StandardErrorPath</key>
    <string>__EARWORM_DIR__/logs/autogen.log</string>
</dict>
</plist>
