# launchd/com.earworm.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Daily episode producer: every day at 07:00, run launchd/daily.sh, which
     recovers crashed runs, tops up the queue only when empty (autogen), then
     drains exactly one topic. The always-on `watch` daemon renders + publishes
     the result -> one new episode per day. The logic lives in the wrapper so it
     composes the stock earworm CLI (no custom flags). Tokens are substituted by
     launchd/install.sh. If the Mac is asleep at 07:00, launchd runs it once on
     the next wake. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.earworm.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__EARWORM_DIR__/launchd/daily.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__EARWORM_DIR__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>7</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__EARWORM_DIR__/logs/daily.log</string>
    <key>StandardErrorPath</key>
    <string>__EARWORM_DIR__/logs/daily.log</string>
</dict>
</plist>
